# rootfs/etc/supervisor.d/supervisord.ini
[supervisord]
user=root
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:mobycron]
command=/usr/bin/mobycron -f /configs/config.json
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:nginx]
command=/docker-entrypoint.sh nginx -g 'daemon off;'
process_name=nginx
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true